# deploy/backend.service
[Unit]
Description=DebtTracker FastAPI Backend
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/debttracker/backend
EnvironmentFile=/var/www/debttracker/.env
ExecStart=/var/www/debttracker/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
